# hlnsrv-timer-1.service
[Unit]
Description=HlnSrv Script Timer 1 Service

[Service]
Type=oneshot
ExecStart=/usr/bin/hlnsrv-script timer_one
