# pkg-arch/hlnsrv-timer-2.service
[Unit]
Description=HlnSrv Script Timer 2 Service

[Service]
Type=oneshot
ExecStart=/usr/bin/hlnsrv-script timer_two
